# DockerFile/supervisord-lamp.conf
[program:mysqld]
command=/mysql-start.sh
numprocs=1
autostart=true
autorestart=false

[program:apache2]
command=/apache2-start.sh
numprocs=1
autostart=true
autorestart=true
